# crontab
*/5 * * * * /app/go-luxpower-timescaledb
